# Quadlet/pc-llm-db.container
[Container]
Mount=type=volume,source=pv-llm-db.volume,target=/data
Image=docker.io/valkey/valkey:9-alpine
Pod=pp-llm.pod
