# {--user}/podman_volume_wireguard.service
# Automatically generated by /usr/lib/systemd/system-generators/podman-system-generator
# 
[X-Volume]
ServiceName=podman_volume_wireguard
VolumeName=wireguard

[Unit]
Wants=network-online.target
After=network-online.target
SourcePath=/etc/containers/systemd/wireguard.volume
RequiresMountsFor=%t/containers

[Service]
ExecStart=/usr/bin/podman volume create --ignore wireguard
SyslogIdentifier=%N
Type=oneshot
RemainAfterExit=yes
