# {--user}/podman_volume_vpn_wireguard.service
# Automatically generated by /usr/lib/systemd/system-generators/podman-system-generator
# 
[X-Volume]
ServiceName=podman_volume_vpn_wireguard
VolumeName=vpn-wireguard

[Unit]
Wants=network-online.target
After=network-online.target
SourcePath=/etc/containers/systemd/vpn-wireguard.volume
RequiresMountsFor=%t/containers

[Service]
ExecStart=/usr/bin/podman volume create --ignore vpn-wireguard
SyslogIdentifier=%N
Type=oneshot
RemainAfterExit=yes
