# {--user}/podman_container_wireguard.service
# Automatically generated by /usr/lib/systemd/system-generators/podman-system-generator
# 
[X-Container]
ServiceName=podman_container_wireguard
ContainerName=wireguard
Image=localhost/wireguard
DropCapability=all
AddCapability=CAP_NET_ADMIN
AutoUpdate=local
Mount=type=volume,source=wireguard.volume,destination=/V
PodmanArgs=--sysctl=net.ipv4.ip_forward=1
PublishPort=51820:51820/udp

[Install]
WantedBy=default.target

[Service]
ExecStartPre=modprobe -r wireguard
ExecStartPre=modprobe wireguard dyndbg=+p
Environment=PODMAN_SYSTEMD_UNIT=%n
KillMode=mixed
ExecStop=/usr/bin/podman rm -v -f -i --cidfile=%t/%N.cid
ExecStopPost=-/usr/bin/podman rm -v -f -i --cidfile=%t/%N.cid
Delegate=yes
Type=notify
NotifyAccess=all
SyslogIdentifier=%N
ExecStart=/usr/bin/podman run --name wireguard --cidfile=%t/%N.cid --replace --rm --cgroups=split --sdnotify=conmon -d --cap-drop all --cap-add cap_net_admin --label io.containers.autoupdate=local --publish 51820:51820/udp --mount type=volume,source=wireguard,destination=/V --sysctl=net.ipv4.ip_forward=1 localhost/wireguard

[Unit]
Wants=network-online.target
After=network-online.target
SourcePath=/etc/containers/systemd/wireguard.container
RequiresMountsFor=%t/containers
Requires=podman_volume_wireguard.service
After=podman_volume_wireguard.service
